# target/dropwizard-deb-package/control/postinst
#!/bin/sh

# SystemD
if which systemctl >/dev/null; then
    echo "Registering SystemD service: scrape"
    systemctl daemon-reload
    systemctl enable "scrape"
    systemctl start "scrape"
    exit 0
fi

# Upstart
if which initctl >/dev/null; then
    echo "Registering Upstart service: scrape"
    initctl reload-configuration
    initctl start "scrape"
    exit 0
fi

# SystemV
if which invoke-rc.d >/dev/null; then
    echo "Registering SystemV service: scrape"
    update-rc.d "scrape" defaults
    invoke-rc.d "scrape" start
    exit 0
fi

echo "Unable to find supported init system."
exit 1
