# target/dropwizard-deb-package/control/prerm
#!/bin/sh

# SystemD
if which systemctl >/dev/null; then
    echo "Removing SystemD service: scrape"
    systemctl stop "scrape"
    systemctl disable "scrape"
    systemctl daemon-reload
    exit 0
fi

# Upstart
if which initctl >/dev/null; then
    echo "Removing Upstart service: scrape"
    initctl stop "scrape"
    initctl reload-configuration
    exit 0
fi

# SystemV
if which invoke-rc.d >/dev/null; then
    echo "Removing SystemV service: scrape"
    invoke-rc.d "scrape" stop
    update-rc.d "scrape" remove
    exit 0
fi

echo "Unable to find supported init system."
exit 1
