# target/dropwizard-deb-package/files/lib/systemd/system/scrape.service
[Unit]
Description=scrape

[Service]
ExecStart=/usr/bin/scrape
Restart=always
OOMScoreAdjust=-1000
User=dropwizard
Group=dropwizard

[Install]
WantedBy=multi-user.target
