# target/dropwizard-deb-package/files/etc/init/scrape.conf
description "scrape"

start on filesystem
stop on runlevel S

respawn
oom never

setuid "dropwizard"
setgid "dropwizard"

exec /usr/bin/scrape
